# src/faultlog/service/faultlog_periodic.service.in
[Unit]
Description=Faultlog application

[Service]
ExecStart=@bindir@/faultlog -r
Type=oneshot
SyslogIdentifier=faultlog

[Install]
#WantedBy=multi-user.target
