# src/faultlog/service/faultlog_bmcboot.service.in
[Unit]
Description=Faultlog application

Wants=mapper-wait@-xyz-openbmc_project-inventory.service
After=mapper-wait@-xyz-openbmc_project-inventory.service
Wants=org.open_power.HardwareIsolation.service
After=org.open_power.HardwareIsolation.service
Wants=mapper-wait@-xyz-openbmc_project-logging.service
After=mapper-wait@-xyz-openbmc_project-logging.service
Wants=com.ibm.VPD.Manager.service
After=com.ibm.VPD.Manager.service

Conflicts=obmc-power-off@0.target
Conflicts=obmc-host-stop@0.target

[Service]
ExecStart=@bindir@/faultlog -r
Type=oneshot
SyslogIdentifier=faultlog

[Install]
WantedBy=multi-user.target
