# src/faultlog/service/faultlog_periodic.service.in
[Unit]
Description=Faultlog application (periodic)
ConditionPathExists=/var/lib/phosphor-software-manager/hostfw/running/DEVTREE
ConditionPathExists=/var/lib/phosphor-software-manager/hostfw/running/GUARD

[Service]
ExecStart=@bindir@/faultlog -r
Type=oneshot
SyslogIdentifier=faultlog

[Install]
#WantedBy=multi-user.target
